# dot_config/systemd/user/backup-s3.service
[Unit]
Description=Backup S3 hebdomadaire (Glacier)

[Service]
Type=oneshot
ExecStart=/home/charles/bin/backup-s3.sh
StandardOutput=append:/home/charles/system-health/backup-s3.log
StandardError=append:/home/charles/system-health/backup-s3.log
Environment=HOME=/home/charles
Environment=PATH=/home/charles/.local/bin:/home/charles/miniconda/bin:/usr/local/bin:/usr/bin:/bin
